# docker/vhosts/supervisor/ssh.conf
[program:ssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true